# ascii-log.cron
#
# Run the ascii-log program at 00:05 every day.
#
30 1 * * * root /usr/bin/ascii-log >/dev/null 2>&1
