# ascii-log.cron
#
# Run the ascii-log program at 01:30 and 13:30 every day.
#
30 1,13 * * * root /usr/bin/ascii-log >/dev/null 2>&1
