# ascii-log.cron
#!/bin/sh

set -e

#
# Run the ascii-log program at 00:05 every day.
#
5 0 * * * /usr/bin/ascii-log >/dev/null 2>&1
