# ascii-log.cron
#
# Run the ascii-log program at 02:30 and 14:30 every day.
#
30 2,14 * * * root /usr/bin/ascii-log >/dev/null 2>&1
